# legacy/hft-optimize.service
[Unit]
Description=HFT System Optimization Service
Documentation=man:systemd.service(5)
Wants=network-online.target
After=network-online.target
Before=multi-user.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/opt/hft/scripts/optimize_system.sh start
ExecStop=/opt/hft/scripts/optimize_system.sh stop
ExecReload=/opt/hft/scripts/optimize_system.sh reload
User=root
Group=root
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=false
PrivateTmp=false
PrivateDevices=false
ProtectSystem=false
ProtectHome=false
ReadWritePaths=/sys /proc /dev

# Resource limits
LimitNOFILE=1048576
LimitNPROC=1048576
LimitMEMLOCK=infinity

# CPU and memory settings
CPUAccounting=true
MemoryAccounting=true

# Restart policy
Restart=no
TimeoutStartSec=300
TimeoutStopSec=60

[Install]
WantedBy=multi-user.target